# configs/clockwork.conf
[program:clockwork]
command=clockwork /root/clock.rb
user=root
pidfile=/clock.pid
stdout_logfile=/plexmediaserver/clock.log
redirect_stderr=true